# autostart/guake.desktop
[Desktop Entry]
Name=Guake
Exec=guake
Type=Application
